# sdlSpeedometer.service
# Systemd resource file
[Unit]
DefaultDependencies=no
Description=Marine instrument application
After=xorg.service

[Service]
Type=simple
EnvironmentFile=/etc/default/sdlSpeedometer
ExecStart=/usr/local/bin/sdlSpeedometer $ARGS
RemainAfterExit=no
Restart=always
RestartSec=1s
User=root

[Install]
WantedBy=sysinit.target
